# etc/systemd/system/statuser.service
[Unit]
Description=Notify to telegram about down services

[Service]
Type=oneshot
ExecStart=/usr/local/bin/send_checks_to_telegram.sh